# ansible/roles/office/files/.config/systemd/user/vdirsyncer.service
[Unit]
Description=Sync vdirsyncer

[Service]
ExecStart=vdirsyncer sync
Type=oneshot
