# ansible/roles/auto_update/files/update.service
[Unit]
Description="A script to update the system"

[Service]
ExecStart=sudo pacman -Syu --noconfirm
